# build/redis/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:redis1]
command=redis-server /usr/local/etc/redis/redis.conf --port 7000 --cluster-config-file nodes-7000.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:redis2]
command=redis-server /usr/local/etc/redis/redis.conf --port 7001 --cluster-config-file nodes-7001.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:redis3]
command=redis-server /usr/local/etc/redis/redis.conf --port 7002 --cluster-config-file nodes-7002.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:redis4]
command=redis-server /usr/local/etc/redis/redis.conf --port 7003 --cluster-config-file nodes-7003.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:redis5]
command=redis-server /usr/local/etc/redis/redis.conf --port 7004 --cluster-config-file nodes-7004.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:redis6]
command=redis-server /usr/local/etc/redis/redis.conf --port 7005 --cluster-config-file nodes-7005.conf
directory=/usr/local/etc/redis
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:setup]
command=/usr/local/bin/setup.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=0
autorestart=false
priority=1
startsecs=5